# docker/root/etc/s6-overlay/s6-rc.d/svc-magnetar/run
#!/usr/bin/with-contenv bash
exec /app/magnetar serve
